# lib/etc/suap-de.conf
[supervisord]
nodaemon=true

[program:cron]
command=/usr/sbin/cron -f
autostart=true
autorestart=true

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true

[program:vscode]
command=/usr/bin/sudo -u suap -i code-server --host 0.0.0.0 --auth password --disable-update-check /app/suap
autostart=true
autorestart=true
